# .config/autostart/gpaste-client.desktop
[Desktop Entry]
Type=Application
Exec=gpaste-client
