# .config/autostart/idle-monitor.desktop
[Desktop Entry]
Type=Application
Exec=idle-monitor
