# .config/autostart/screensaver.desktop
[Desktop Entry]
Type=Application
Exec=screensaver
